# solarupdate/crontab
### * * * * * root echo "hello alive" >> /dvol/log/solarupdate.log 2>&1
*/5 6-19 * * * root python3 /root/solar_update.py >> /dvol/log/solarupdate.log 2>&1
